# etc/yservice/supervisor.conf
[supervisord]
nodaemon = true

[program:nginx]
command = service nginx restart

[program:uwsgi]
directory = /etc/yservice
command = uwsgi -i yservice.ini
