# docker/supervisor.conf
[supervisord]
nodaemon = true

[program:nginx]
command = service nginx restart

[program:uwsgi]
directory = /home/opnfv/repos/yardstick-service/docker
command = uwsgi -i yardstick-service.ini
