# docker/supervisor.conf
[supervisord]
nodaemon = true

[program:nginx]
command = service nginx restart

[program:uwsgi]
directory = /home/publiccloud/docker
command = uwsgi -i publiccloud.ini
